# systemd.user/clean-vim-undo.timer
[Unit]
Description=Clean VIM undofiles

[Timer]
OnBootSec=10s

[Install]
WantedBy=timers.target
